# wifi-bridge-listener.service
[Unit]
Description=WiFi Bridge Listener Service (NetworkManager)
# Wait until the network stack is fully configured and online
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# User MUST be root for nmcli
User=root
# ExecStart: Use the absolute path in the system location
ExecStart=/usr/bin/python3 /usr/local/sbin/RPZero2WListener.py

# Restart the service if it fails
Restart=on-failure
RestartSec=5

[Install]
# Start the service during the multi-user boot target
WantedBy=multi-user.target

